# virtual-boxes/ubuntu-ff/selenium-node.conf
# selenium-node

start on runlevel [2345]
stop on runlevel [016]
respawn
exec java -jar selenium-server-standalone-2.46.0.jar